# ozer-garant-bot.service
[Unit]
Description=OZER Garant Telegram Bot
After=network.target mysql.service

[Service]
Type=simple
User=www-data
WorkingDirectory=/path/to/your/bot
Environment=PATH=/usr/bin:/usr/local/bin
ExecStart=/usr/bin/python3 /path/to/your/bot/bot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target